# etc/udev/rules.d/51-fanout.rules
ACTION=="add", SUBSYSTEM=="module", DEVPATH=="/module/fanout", RUN+="/opt/em/bin/create-fanout-devs.sh"
