# etc/systemd/system/fake-hwclock.service
[Unit]
Description=Fake Hardware Clock Service

[Service]
Type=forking
ExecStart=/usr/bin/fake-hwclock
Nice=-5
Restart=always
