# etc/systemd/system/sensors.timer
[Unit]
Description=Log LM sensors output every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=60min
