# etc/systemd/system/sensors.service
[Unit]
Description=LM Sensors

[Service]
Type=idle
ExecStart=/usr/bin/sensors
