# etc/systemd/system/serial-setup.service
[Unit]
Description=Set serial port options
After=systemd-udev-trigger.service

[Service]
Type=oneshot
ExecStart=/opt/em/bin/serial-setup.sh
