# etc/systemd/system/smartctl.timer
[Unit]
Description=Log SMART health check results every 8 hours

[Timer]
OnBootSec=1min
OnUnitActiveSec=8h
