# etc/systemd/system/smartctl.service
[Unit]
Description=SMART health check

[Service]
Type=oneshot
ExecStart=-/usr/sbin/smartctl -x /dev/sda
ExecStart=-/usr/sbin/smartctl -x /dev/sdb
